# docs/assets/code/systemd/home-bkp-nas.timer
[Unit]
Description=Run NAS Backup Weekdays at 7PM

[Timer]
# Mon-Fri at 19:00 (7 PM)
OnCalendar=Mon..Fri 19:00:00
# If the computer is off/asleep at 6pm, run immediately when it wakes up
Persistent=true
# Randomize delay by 5 mins to prevent exact-second spikes (optional but good practice)
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target