# docs/assets/code/systemd/mnt-storage_server-NVMe.automount
[Unit]  
Description=Automount for ZimaOS Quick-Storage  
  
[Automount]  
Where=/mnt/storage_server/NVMe  
TimeoutIdleSec=600  
  
[Install]  
WantedBy=multi-user.target