# docs/assets/code/systemd/mnt-storage_server-Quick_Storage.automount
[Unit]  
Description=Automount for ZimaOS NAS Quick-Storage  
  
[Automount]  
Where=/mnt/storage_server/Quick_Storage  
# Unmount if idle for 10 minutes (optional, saves resources)  
TimeoutIdleSec=600  
  
[Install]  
WantedBy=multi-user.target